# marketing/automation/com.longhu.obsidian-ai-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.longhu.obsidian-ai-news</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/longhuadmin/obsidian-inbox/scripts/run_ai_news_daily.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/longhuadmin/obsidian-inbox/logs/ai-news-daily.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/longhuadmin/obsidian-inbox/logs/ai-news-daily.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/longhuadmin/obsidian-inbox</string>
</dict>
</plist>
